# systemd/.config/systemd/user/devilspie2.service
[Unit]
Description:devilspie

[Service]
ExecStart=/usr/bin/devilspie2
